# cron/crontab
* * * * * echo "Crontab is working - watchdog 1"